# playbook/project/roles/micado_master/build/files/iptables/ip6tables.service
[Unit]
Description=Restore iptables firewall rules
After=iptables.service

[Service]
Type=oneshot
ExecStart=/sbin/ip6tables-restore -n /etc/ip6tables.conf

[Install]
WantedBy=multi-user.target

